# install/inits/systemd/system/tyk-gateway-lua.service
[Unit]
Description=Tyk API Gateway (LUA Support)
 
[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/tyk-gateway/tyk-lua --conf=/opt/tyk-gateway/tyk.conf
Restart=always
WorkingDirectory=/opt/tyk-gateway

[Install]
WantedBy=multi-user.target
