# install/inits/systemd/system/tyk-gateway.service
[Unit]
Description=Tyk API Gateway
 
[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/tyk-gateway/tyk --conf=/opt/tyk-gateway/tyk.conf
Restart=always

[Install]
WantedBy=multi-user.target
