# terraform/shared/rc.local
sysctl fs.file-max=80000
sysctl vm.overcommit_memory=1
sysctl net.core.somaxconn=1024
echo never | tee /sys/kernel/mm/transparent_hugepage/enabled
